# templates/default/etc/cron.d/housekeeping.erb
27 1 * * * root find /var/lib/docker/volumes/transmission-config/_data/torrents -mindepth 1 -maxdepth 1 -mtime +7 -type f -delete && systemctl restart docker-transmission.service
37 * * * * root if [ "$(df / | tail -n1 | tr -s " " | cut -f5 -d" " | cut -f1 -d\%)" -gt "70" ]; then find /home/media/downloads -mindepth 2 -mtime +30 -delete; fi
47 * * * * root if [ "$(df / | tail -n1 | tr -s " " | cut -f5 -d" " | cut -f1 -d\%)" -gt "80" ]; then find /home/media/downloads -mindepth 2 -mtime +14 -delete; fi
57 * * * * root if [ "$(df / | tail -n1 | tr -s " " | cut -f5 -d" " | cut -f1 -d\%)" -gt "90" ]; then find /home/media/downloads -mindepth 2 -mtime +7 -delete; fi
