# templates/default/etc/systemd/system/docker-sickgear.service.erb
#
# Docker + SickGear systemd service
#
[Unit]
Description=SickGear Docker Container
After=network.target docker.socket
Requires=docker.socket

[Service]
RestartSec=10
Restart=always

Environment="NAME=sickgear"
Environment="IMG=deed02392/sickgear"

# Clean-up bad state if still hanging around
ExecStartPre=-/usr/bin/docker rm -f $NAME

# Attempt to pull new image for security updates
ExecStartPre=-/usr/bin/docker pull $IMG

# Main process
ExecStart=/usr/bin/docker run --rm --name ${NAME} \
  -v /home/media/downloads:/incoming \
  -v /home/media/shows:/tv \
  -v /home/media/sickgear:/data \
  -e APP_UID=1001 \
  -e APP_GID=1001 \
  -e TZ=UTC \
  -p 172.17.0.1:8081:8081 \
  ${IMG}

[Install]
WantedBy=multi-user.target
