# templates/default/etc/systemd/system/iptables-restore.service.erb
[Unit]
Description=Restore iptables firewall rules
Before=network.target
Conflicts=shutdown.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c "/sbin/ipset restore < /etc/ipset.state"
ExecStart=/usr/sbin/iptables-restore /etc/iptables.rules
ExecStart=/usr/sbin/ip6tables-restore /etc/ip6tables.rules

[Install]
WantedBy=basic.target
